# src/system/etc/systemd/system/pi-boot-network-base.service
[Unit]
Description=PI5: base routing & policy rules before NordVPN wrapper
Requires=systemd-sysctl.service nftables.service wg0.service
After=systemd-sysctl.service nftables.service wg0.service network-online.target
Wants=network-online.target
Before=nordvpn-wrapper.service

[Service]
Type=oneshot
# kleine Guard-Waits, damit wir nicht ins Leere greifen:
TimeoutStartSec=30
ExecStartPre=/bin/sh -c 'until nft -c list ruleset >/dev/null 2>&1; do echo "[BOOT] waiting for nftables…"; sleep 0.5; done'
# ExecStartPre=/bin/sh -c 'until ip -4 addr show wg0 | grep -q "10\\.6\\."; do echo "[BOOT] waiting for wg0 address…"; sleep 0.5; done'
ExecStart=/usr/local/sbin/pi-boot-network-base.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
